# modules/daemontools/files/skipjack-svscan.upstart.conf
# svscan - daemontools - skipjack

start on filesystem
stop on runlevel [06]

respawn
exec /usr/bin/svscan.skipjack
