# modules/devil_ray/files/devil_ray.conf
# devil_ray - S3 proxy for requester-pays buckets

description "Devil Ray"

start on started network-interface
stop on stopping network-interface

chdir /opt/devil_ray
exec ./devil_ray.py
